# keepers/deploy/fusd-alert-webhook.service
[Unit]
Description=fUSD alert webhook bridge (monitor criticals + liveness -> Discord)

[Service]
Type=oneshot
User=@RUNUSER@
EnvironmentFile=/etc/fusd/keeper.env
ExecStart=@REPO@/keepers/alert-webhook.sh
# Sandbox. Pure bash/curl/jq (no node), but ProtectHome=read-only (NOT tmpfs) because @REPO@
# lives under the operator's home per the README. Only the state dir is writable.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
ReadWritePaths=/var/lib/fusd
InaccessiblePaths=-@HOME@/.config/solana
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectKernelLogs=true
ProtectClock=true
ProtectHostname=true
RestrictSUIDSGID=true
RestrictRealtime=true
LockPersonality=true
